# Code/RaspberryPi5/SpeakerConfiguration/rc.local
#!/bin/bash
# Restore ALSA settings
sudo alsactl restore

exit 0
